# audio/audio-fix.service
[Unit]
Description=Audio Fix

[Service]
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
