# audio/audio-fix.service
[Unit]
Description=Audio Fix

[Service]
ExecStart=/usr/local/bin/necessary-verbs.sh
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
